# roles/collectd/templates/supervisor.conf
[program:collectd]
command=/home/collectd/collectd-web/runserver.py 127.0.0.1 8093
directory=/home/collectd/collectd-web
user=collectd
numprocs=1
autostart=true
autorestart=true
startsecs=10
